# chart/files/etc/profile.d/alias.sh
#!/bin/bash

alias sync-curriculum="gitpuller /nh/curriculum main ~/curriculum || gitpuller https://github.com/neurohackademy/nh2021-curriculum.git master ~/curriculum || true"
